# luci-app-oaf/root/etc/uci-defaults/91_luci-oaf-2
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	del ucitrack.@appfilter[-1]
	commit ucitrack
EOF

# remove LuCI cache
rm -rf /tmp/luci-indexcache /tmp/luci-modulecache
